# spushnik-server/src/deb/control/postinst
#!/bin/bash
set -e

invoke-rc.d spushnik start	 3>1

#DEBHELPER#

exit 0

