# spushnik-server/src/deb/control/postrm
#!/bin/bash
set -e

if [ -f /etc/init.d/spushnik ]; then
	invoke-rc.d spushnik stop 2>&1
fi 	

exit 0

