# spushnik-server/src/deb/control/postrm
#!/bin/bash
set -e

invoke-rc.d spushnik stop 2>&1 	

exit 0

