# push-server/src/deb/control/prerm
#!/bin/sh -e

if [ -f /etc/init.d/opush ]; then
	invoke-rc.d opush stop 3>1
fi